# recipes-extended/automountd/files/automount.service
[Unit]
Description=Automounter using udevadm

[Service]
Type=simple
ExecStart=/etc/udev/scripts/automountd
Restart=on-failure

